# mtk-linkit/files/etc/uci-defaults/54_linkit_version
#!/bin/sh

uci set system.@system[0].linkit_firstboot=1
uci commit system
